# snapshot/mail-snapshot-monthly.service
[Unit]
Description=Snapshot mail and copy to archive as monthly snapshot

[Service]
ExecStart=/usr/local/bin/mail-snapshot.sh monthly
