# snapshot/mail-snapshot-daily.service
[Unit]
Description=Snapshot mail and copy to archive as daily snapshot

[Service]
ExecStart=/usr/local/bin/mail-snapshot.sh daily
