# snapshot/mail-snapshot-weekly.timer
[Unit]
Description=Snapshot mail weekly

[Timer]
OnCalendar=weekly
Unit=mail-snapshot-weekly.service

[Install]
WantedBy=timers.target
