# snapshot/mail-snapshot-monthly.timer
[Unit]
Description=Snapshot mail monthly

[Timer]
OnCalendar=monthly
Unit=mail-snapshot-monthly.service

[Install]
WantedBy=timers.target
