# snapshot/mail-snapshot-daily.timer
[Unit]
Description=Snapshot mail daily

[Timer]
OnCalendar=daily
Unit=mail-snapshot-daily.service

[Install]
WantedBy=timers.target
